# roles/init_and_unseal/templates/vault-healthcheck.service.j2
[Unit]
Description=Vault Seal Status Health Check
Documentation=https://www.vaultproject.io/docs/
# No After/BindsTo on a "vault.service" unit: Vault runs as a Docker
# container, not a native systemd service -- that unit never exists.
# BindsTo on a nonexistent unit made every timer-triggered start of this
# service fail outright with "Dependency failed", so the healthcheck (and
# therefore its own systemctl start vault-auto-unseal.service call) never
# actually ran at all.
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vault-healthcheck.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vault-healthcheck
User=root
Group=root
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
