# roles/init_and_unseal/templates/vault-auto-unseal.service.j2
[Unit]
Description=Vault Auto-Unseal Service
Documentation=https://www.vaultproject.io/docs/
# No Requires/After on a "vault.service" unit: Vault runs as a Docker
# container (managed by community.docker.docker_container), not a native
# systemd service -- that unit never exists, and a dependency on it
# previously made every attempt to start this service fail outright with
# "Dependency failed" / result 'dependency', silently preventing
# auto-unseal from ever running at all.
After=network-online.target docker.service
ConditionPathExists={{ vault_unseal_script_path }}

[Service]
Type=oneshot
ExecStart={{ vault_unseal_script_path }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vault-auto-unseal
User=root
Group=root
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes

# P1-8: Retry logic in systemd (in addition to script retry loop)
# Systemd will retry this service if it fails
Restart=on-failure
RestartSec=10
StartLimitInterval=120
StartLimitBurst=3

RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
